# trunk/bilibot/bilibot_bringup/upstart/bilibot.conf
description	"bringup bilibot"

start on net-device-up IFACE=wlan0
stop on net-device-down IFACE=wlan0

console output
respawn

exec bilibot-start
pre-stop exec bilibot-stop
